# deploy/launchd/com.tunnelmesh.peer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tunnelmesh.peer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/tunnelmesh</string>
        <string>join</string>
        <string>--config</string>
        <string>/etc/tunnelmesh/peer.yaml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/var/log/tunnelmesh/peer.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/tunnelmesh/peer.error.log</string>

    <key>WorkingDirectory</key>
    <string>/var/lib/tunnelmesh</string>

    <key>UserName</key>
    <string>root</string>
</dict>
</plist>
